# supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700
chown = root:root

[supervisord]
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:myapp]
command = /usr/local/bin/daphne -b 0.0.0.0 -p 8000 djangoProject.asgi:application
directory = /root/project/djangoProject
user = root
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/django.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5